# roles/bitcoind/templates/upstart.conf
## This file is autogenerated by ansible-honeybadger. Do not modify it directly.
## Instead, update the template in ansible-honeybadger, and redeploy. Maybe
## submit a pull request?

description "bitcoind"

start on filesystem
stop on runlevel [!2345]

expect fork

setuid bitcoind
setgid bitcoind

nice {{ bitcoind_nice }}

# we don't set respawn here: monit does a better job of handling restarts

{% if bitcoind_ionice %}
pre-start script
    ionice -c{{bitcoind_ionice_class}} -n{{bitcoind_ionice_class_data}} -p $PPID
end script
{% endif %}

{% if bitcoind_cpulimit %}
post-start script
    sleep 3 # wait for the pid file to get written out
    PID=$(cat /var/lib/bitcoind/bitcoind.pid)
    /opt/cpulimit/src/cpulimit -il {{ bitcoind_cpulimit }} -p $PID &
end script
{% endif %}

script
    exec /usr/bin/bitcoind -conf=/etc/bitcoin.conf \
                           -datadir=/var/lib/bitcoind \
                           -daemon
end script
